# nonweb/pkg/SystemHealth/0.1/syscheck.cron
#! /bin/sh

subject="OLAC ALERT ($(date +'%b %d %H:%M %Y'))"
recipients="$(olacvar tech_email) $(olacvar olac_admin_email | sed 's/,/ /')"
$(dirname $0)/syscheck enabled-daily | grep -v 'PASS$' | \
    /usr/bin/mail -e -s "$subject" $recipients
